# supervisor/php-fpm_supervisor.conf

[program:php-fpm]
command=/usr/local/php/sbin/php-fpm --nodaemonize -R



